# supervisor-app.conf
[supervisord]
nodaemon=true


[program:meta_rest_app]
command=/usr/bin/python3.8 manage.py runserver
autostart=true
autorestart=true
redirect_stderr=true
