# supervisor-app.conf
[supervisord]
nodaemon=true


[program:meta_rest_app]
command=python manage.py
startsecs = 0
autorestart = false
http://supervisord.org/configuration.html

startsecs
