# supervisor-app.conf
[program:meta_rest_app]
command=/home/ubuntu/python/run_local_server.sh
autostart=true
autorestart=true
redirect_stderr=true
