# supervisor-app.conf
[supervisord]
nodaemon=true


[program:meta_rest_app]
command=/home/ubuntu/python/virtualenvs/python manage.py runserver
autostart=true
autorestart=true
redirect_stderr=true
